# BlueOnyx/5106R/ui/base-console.mod/glue/cron.d/pam_abl_purge.cron
# base-console.mod - Cronjob that purges the pam_abl databases.
*/15 * * * * root /usr/bin/pam_abl --purge > /dev/null 2>&1 &
