# BlueOnyx/5207R/ui/base-ssl.mod/src/blueonyx-le-acme/acme/crontab.cron
30 0 * * * /usr/sausalito/acme/acme.sh --cron --home /usr/sausalito/acme --config-home /usr/sausalito/acme/data > /dev/null
