# BlueOnyx/5107R/swatch/swatch.cron
# swatch is the system watcher
*/15 * * * *    root    /usr/sausalito/sbin/swatch.sh
