# BlueOnyx/5107R/common/base-blueonyx.mod/glue/etc/cron.d/db_checkpoint.cron
# exec db_checkpoint.pl
03,18,33,48 * * * *	root	/usr/sausalito/sbin/db_checkpoint.pl
