# BlueQuartz/5100WG/tags/OSS_1_4/ui/base-quotastats.mod/glue/etc/cron.d/quotastats.cron
# Run quota stats every 6 hours
# (amdisk.cron runs it at midnight with the email option)
0 6,12,18 * * *    root    /usr/sausalito/sbin/get_quota_stats.pl >/dev/null 2>/dev/null
