# BlueOnyx/tmp/base-apache.diff/base-apache.mod/glue/etc/cron.d/suphp_session_file_clear.cron
*/5 * * * * root    /usr/local/sbin/suphp_session_file_clear.sh > /dev/null 2>&1
