# BlueOnyx/5210R/ui/base-console.mod/glue/cron.daily/pam_abl.sh
#!/bin/sh
# Purge all PAM_ABL records and databases:
rm -f /var/lib/pam_abl/*

