# BlueOnyx/5209R/platform/base-apache.mod/glue/php-fpm/php-fpm-7.9.service
[Unit]
Description=The PHP-7.9 FastCGI Process Manager
After=syslog.target network.target

[Service]
Type=notify
PIDFile=/run/php-fpm/php-fpm-7.9.pid
EnvironmentFile=/etc/sysconfig/php-fpm
ExecStart=/home/solarspeed/php-7.9/sbin/php-fpm --nodaemonize
ExecReload=/bin/kill -USR2 $MAINPID
PrivateTmp=true
Restart=on-failure

[Install]
WantedBy=multi-user.target

