# BlueOnyx/5209R/utils/swatch/swatch.cron
# swatch is the system watcher
*/15 * * * *    root    /usr/sausalito/sbin/swatch.sh >/dev/null 2>&1