# BlueQuartz/5100WG/tags/OSS_1_4/ui/base-disk.mod/glue/etc/cron.d/amdisk.cron
# Check the disk quotas and send email notices if necessary
# Do this every night at 11:30 PM
30 23 * * *	root	/usr/sausalito/sbin/get_quota_stats.pl -m >/dev/null 2>/dev/null
