# BlueOnyx/5106R/ui/base-email.mod/glue/etc/cron.d/break-smrsh-loop.cron
*/5 * * * *	root	/usr/local/sbin/break-smrsh-loop.pl > /dev/null 2>&1
