# BlueQuartz/5100WG/branches/dev/ui/base-webstats.mod/glue/etc/cron.d/webstats.cron
# update the web stats every hour
0 0-23 * * *    root    /usr/sausalito/sbin/parseWebLog.pl /var/log/httpd/access >/dev/null 2>/dev/null
