# BlueQuartz/5100R/trunk/swatch/swatch.cron
# swatch is the system watcher
*/15 * * * *    root    /usr/sbin/swatch -c /etc/swatch.conf
