# BlueOnyx/5207R/ui/base-memcache.mod/glue/cron.weekly/restart_cce.cron
#!/bin/bash
/etc/init.d/cced.init restart
